# SOURCES/nagrestconf-1/debian/postinst
#!/bin/sh
# postinst script for thruk
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        # enable modules we need
        a2enmod php5
        a2enmod auth_basic
        a2enmod rewrite

        usermod -a -G nagios www-data
        #usermod -a -G nagiocmd www-data

        # reload new apache config

        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d apache2 reload
        else
            /etc/init.d/apache2 reload
        fi

        # Add crontab

        if test -w /var/spool/cron/crontab/root; then
            # Delete crontab first
            crontab -l | sed '/nagios_restart_request/d' \
              | crontab -
        fi
        echo "* * * * * /usr/bin/test -e /tmp/nagios_restart_request && ( /bin/rm /tmp/nagios_restart_request; /usr/bin/restart_nagios; )" | crontab -

        # Add sudoers entry

        sed -i 's/^\(Defaults[[:space:]]*requiretty\)/#\1/' /etc/sudoers
        sed -i '/\/usr\/bin\/csv2nag -y all/d' /etc/sudoers
        cat >>/etc/sudoers <<EnD
%nagios ALL = NOPASSWD: /usr/sbin/nagios3 -v *, /usr/bin/csv2nag -y all

EnD

        # Add the line slc_configure wants

        if ! grep -qs "<SERVICE_LINE_CFG_ENTRY>" /etc/nagios3/nagios.cfg; then
            sed -i '/SERVICE_LINE_CFG_ENTRY/d' /etc/nagios3/nagios.cfg
            echo "## Next line added by nagrestconf"  >>/etc/nagios3/nagios.cfg
            echo "<SERVICE_LINE_CFG_ENTRY>" >>/etc/nagios3/nagios.cfg

            # Comment out cfg_ lines in nagios.cfg

            sed -i \
                's/^[[:space:]]*cfg_/# - commented out by nagrestconf - cfg_/' \
                /etc/nagios3/nagios.cfg

            slc_configure --folder=local
        fi

        echo "Nagrestconf has been configured for http://127.0.0.1/nagrestconf/."
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

